# manage/topics/redis/templates/redis.ini
[program:redis]
{% if use_conf_file %}
command=redis-server /home/pfsc/redis.conf
{% else %}
command=redis-server
{% endif %}
priority=100
user=pfsc
autostart=true
