# manage/topics/redis/templates/redisgraph.ini
[program:redisgraph]
{% if use_conf_file %}
command=redis-server /home/pfsc/redis.conf
{% else %}
command=redis-server --dir /proofscape/graphdb/re --loadmodule /usr/lib/redis/modules/redisgraph.so
{% endif %}
priority=100
user=pfsc
autostart=true
